# internal/cli/schedule/templates/sync.timer.tmpl
[Unit]
Description=prosa sync timer

[Timer]
OnBootSec=1min
OnUnitActiveSec={{ .IntervalSpec }}
Persistent=true

[Install]
WantedBy=timers.target
