# internal/cli/schedule/templates/sync.service.tmpl
[Unit]
Description=prosa sync — import AI agent session histories

[Service]
Type=oneshot
ExecStart={{ .Binary }} sync
